# .config/supervisord.conf
[supervisord]
logfile = /var/log/supervisor.log
user = root

[program:gunicorn]
command = gunicorn -c /srv/Todoapp/.config/gunicorn.py config.wsgi

[program:nginx]
command = nginx -g "daemon off;"